# vps-setup/animation.service
[Unit]
Description=Animation Server for LiveStream
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/liveStream/animation-server
ExecStart=/root/.nvm/versions/node/v18.20.8/bin/node /home/liveStream/animation-server/server.js
Restart=on-failure
RestartSec=10
Environment=ANIMATION_PORT=3003

[Install]
WantedBy=multi-user.target
